# Desktop/daemon/linux/predictmv.service
[Unit]
Description=PredictMV Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=predict-mv
Group=predict-mv
WorkingDirectory=/opt/predict-mv/daemon
Environment=XDG_CONFIG_HOME=/var/lib/predict-mv/config
ExecStart=/opt/predict-mv/daemon/PredictMV
Restart=always
RestartSec=5
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
